# .utils/systemd_units/hibernate-on-low-battery.service
[Unit]
Description=Hibernate the system when battery drops below critical threshold
Documentation=man:systemd-sleep

[Service]
Type=simple
ExecStart=/usr/local/sbin/hibernate-on-low-battery.sh

[Install]
WantedBy=multi-user.target
